# deploy/incognito-followup.service
[Unit]
Description=Incognito follow-up runner

[Service]
Type=oneshot
# INCOGNITO_PASSWORD must be set in the container environment for non-interactive use.
# Pass it via podman run -e or set Environment= here if not using container secrets.
ExecStart=podman exec incognito incognito follow-up --auto
