# deploy/incognito-followup.service
[Unit]
Description=Incognito follow-up runner

[Service]
Type=oneshot
ExecStart=podman exec incognito incognito follow-up --auto
