# deploy/incognito-rescan.service
[Unit]
Description=Incognito re-scan runner

[Service]
Type=oneshot
ExecStart=podman exec incognito incognito rescan
