# deploy/incognito-followup.timer
[Unit]
Description=Incognito daily follow-up check

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
